# misc/tmmvision.service
[Service]
ExecStart=/home/git/steamworks-vision/start.sh

[Install]
WantedBy=default.target
